# certificate-manager.service
[Unit]
Description=Certificate Manager - Gestionnaire de certificats X.509
Documentation=https://github.com/your-org/certificate-manager
After=network.target
Wants=network.target

[Service]
Type=simple
User=certmanager
Group=certmanager
WorkingDirectory=/opt/certificate-manager
ExecStart=/opt/certificate-manager/venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8080
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Variables d'environnement
Environment=CERTS_DIR=/data/certs
Environment=PORT=8080
Environment=HOST=0.0.0.0
Environment=SCAN_INTERVAL_SEC=300
Environment=THRESHOLD_DAYS=30
Environment=SECRET_KEY=change-this-in-production
# Décommentez et configurez pour l'authentification
# Environment=AUTH_TOKEN=your-secret-token-here
# Environment=BASIC_USER=admin
# Environment=BASIC_PASS=your-password-here

# Sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/certificate-manager/data
ReadOnlyPaths=/data/certs

# Limites de ressources
LimitNOFILE=65536
MemoryMax=512M

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=certificate-manager

[Install]
WantedBy=multi-user.target
